# etc/cron.d/projector
# The pi won't be ready when @reboot triggers, so wait before we start
@reboot			root sleep 15 && /root/startup-slide
